# packages/relay/config/river-relay-reset.service
[Unit]
Description=Reset River Relay systemd failure counter

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl reset-failed river-relay.service
